# debian/couchdb.upstart
# couchdb - a RESTful document oriented database

description "Start Apache CouchDB"
author "CouchDB Developers <dev@couchdb.apache.org>"

start on filesystem and static-network-up
stop on deconfiguring-networking
respawn

exec su couchdb -c /opt/couchdb/bin/couchdb

